# sawi/templates/django_supervisor.conf
[program:oep]
command = /srv/oep/bin/start.sh
stdout_logfile = /srv/oep/log/supervisor.log
redirect_stderr = true
environment = LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8